# rom/etc/hotplug.d/button/53-wps-repacd
#
# Copyright (c) 2015 Qualcomm Atheros, Inc..
#
# All Rights Reserved.
# Qualcomm Atheros Confidential and Proprietary.
#

if [ "$ACTION" = "released" -a "$BUTTON" = "wps" ]; then

fi
