# tmp/mnt/circle/shares/usr/bin/service/sitedbd/run
#!/bin/sh
CIRCLE_ROOT=`cat /tmp/CIRCLE_ROOT`
exec $CIRCLE_ROOT/sitedbd
