# tmp/mnt/circle/shares/usr/bin/service/webd/run
#!/bin/sh
CIRCLE_ROOT=`cat /tmp/CIRCLE_ROOT`
exec $CIRCLE_ROOT/webd -u root -E $CIRCLE_ROOT/mini_httpd.pem  -D
