# rom/usr/lib/opkg/info/armor.list
/usr/share/armor/change_cloud_server.sh
/usr/share/armor/BD_START.sh
/usr/share/armor/bd_account
/lib/libbdbroker.so
/usr/share/armor/bdbroker-client-thread
/usr/share/armor/bdbroker-client-loop
/usr/share/armor/upgrade_bd_cfg.sh
/usr/share/armor/bdupd_start_schedule.sh
/usr/share/armor/bdbroker-client-cmd
/etc/init.d/armor_init
/usr/share/armor/get_armor_status
/usr/share/armor/bdupd_start.sh
/lib/armor/phase2-upd.tar.gz
