# rom/usr/lib/opkg/info/control-firmware.list
/usr/sbin/control-firmware
/etc/init.d/control_firmware
/usr/sbin/boot-notification
