# rom/usr/lib/opkg/info/detcable.list
/etc/init.d/detcable
/usr/sbin/port_detect.sh
/usr/bin/detcable
/usr/bin/flushRoute
