# rom/usr/lib/opkg/info/dni-network.list
/etc/init.d/network.init
/etc/init.d/net-wan
/etc/sw_extender.conf
/etc/sw.conf
/etc/init.d/net-br
/etc/init.d/net-lan
/etc/init.d/net-lan-ext
/etc/init.d/net-br-dhcpc-helper
/etc/config/network
/etc/init.d/switch.init
