# rom/usr/lib/opkg/info/dni-ntpclient-ext.list
/usr/sbin/ntpclient-ext
/etc/init.d/ntpclient-ext
