# rom/usr/lib/opkg/info/dnibusybox.list
/bin/dmesg
/usr/bin/id
/bin/login
/sbin/blockdev
/bin/rm
/usr/bin/wget
/usr/bin/sha512sum
/sbin/lsmod
/bin/sh
/usr/bin/lsof
/sbin/vconfig
/bin/su
/usr/bin/test
/sbin/mkdosfs
/sbin/depmod
/usr/bin/truncate
/usr/bin/unlzop
/bin/zcat
/bin/fgrep
/usr/sbin/killall5
/bin/ash
/sbin/halt
/sbin/reboot
/bin/true
/sbin/route
/usr/bin/md5sum
/usr/bin/free
/bin/fatattr
/usr/bin/time
/usr/bin/passwd
/sbin/klogd
/usr/bin/traceroute
/usr/sbin/rdev
/sbin/acpid
/bin/hostname
/bin/vi
/sbin/md
/usr/bin/tftp
/bin/pwd
/sbin/modprobe
/sbin/mm
/usr/bin/nc
/sbin/mke2fs
/bin/chown
/etc/init.d/cron
/usr/sbin/rdate
/usr/sbin/addgroup
/bin/mount
/usr/bin/xzcat
/usr/sbin/i2cget
/bin/sync
/usr/bin/shuf
/usr/sbin/adduser
/usr/bin/dirname
/usr/bin/bzip2
/sbin/rmmod
/usr/bin/hostid
/bin/bash
/sbin/hwclock
/usr/bin/ttysize
/usr/bin/traceroute6
/sbin/insmod
/usr/bin/tac
/bin/ping6
/etc/bandcheck/band-check
/usr/bin/basename
/sbin/slattach
/bin/dnsdomainname
/usr/sbin/delgroup
/bin/echo
/usr/bin/less
/usr/bin/dos2unix
/usr/bin/unix2dos
/usr/sbin/chroot
/usr/bin/unexpand
/sbin/devmem
/usr/sbin/crond
/usr/bin/head
/bin/gunzip
/bin/netstat
/usr/bin/fgconsole
/usr/sbin/ftpd
/bin/busybox
/usr/bin/pwdx
/bin/grep
/usr/bin/clear
/usr/bin/seq
/bin/sleep
/bin/nice
/usr/bin/diff
/sbin/watchdog
/usr/bin/groups
/usr/bin/sha1sum
/bin/kill
/bin/mktemp
/usr/bin/tr
/usr/bin/whois
/usr/bin/mesg
/bin/mknod
/bin/rmdir
/usr/bin/uniq
/usr/bin/tee
/bin/getopt
/usr/bin/flock
/usr/sbin/powertop
/usr/bin/bzcat
/usr/bin/hexdump
/usr/bin/split
/usr/bin/sha256sum
/bin/date
/usr/sbin/dsyslog
/usr/bin/wc
/usr/bin/nslookup
/usr/bin/[
/bin/base64
/usr/bin/pgrep
/bin/umount
/usr/bin/tail
/bin/mpstat
/bin/stat
/usr/bin/pkill
/usr/bin/mkfifo
/usr/sbin/i2cdump
/usr/sbin/i2cdetect
/bin/chgrp
/usr/bin/xz
/usr/bin/pscan
/bin/cp
/usr/bin/mkpasswd
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/lock
/usr/bin/expand
/bin/df
/bin/false
/usr/bin/sha3sum
/usr/bin/killall
/usr/bin/sort
/bin/mkdir
/sbin/ifconfig
/usr/bin/reset
/usr/bin/fuser
/sbin/start-stop-daemon
/bin/ionice
/usr/bin/pmap
/usr/bin/chrt
/usr/bin/microcom
/sbin/fbsplash
/bin/usleep
/bin/cat
/etc/init.d/syslogd
/usr/sbin/arping
/bin/lzop
/sbin/fstrim
/usr/bin/tcpsvd
/sbin/mkfs.vfat
/usr/bin/smemcap
/usr/bin/telnet
/usr/sbin/tftpd
/sbin/getty
/bin/touch
/bin/tar
/bin/rev
/usr/bin/lsusb
/bin/iostat
/usr/bin/yes
/bin/conspy
/usr/bin/top
/sbin/fdisk
/usr/bin/expr
/sbin/switch_root
/usr/bin/env
/usr/bin/showkey
/sbin/mkfs.ext2
/sbin/sysctl
/sbin/syslogd
/bin/sed
/usr/bin/logger
/sbin/uevent
/bin/chmod
/bin/setserial
/usr/bin/awk
/bin/uname
/usr/bin/which
/sbin/modinfo
/sbin/zcip
/usr/bin/find
/bin/ln
/usr/bin/users
/usr/bin/pstree
/usr/sbin/chat
/bin/ls
/usr/bin/cryptpw
/sbin/pivot_root
/sbin/findfs
/bin/kbd_mode
/usr/bin/cut
/usr/bin/du
/bin/mv
/usr/sbin/readahead
/usr/sbin/telnetd
/usr/bin/crontab
/bin/cttyhack
/usr/bin/printf
/usr/bin/unxz
/usr/sbin/cmd_cron
/usr/bin/lspci
/bin/ping
/usr/bin/timeout
/usr/bin/bunzip2
/usr/sbin/chpasswd
/usr/bin/uptime
/bin/fsync
/sbin/tunctl
/usr/sbin/deluser
/bin/pidof
/bin/egrep
/usr/bin/readlink
/usr/bin/strings
/sbin/arp
/usr/bin/unlink
/usr/bin/udpsvd
/bin/ps
/usr/bin/lzopcat
/usr/sbin/i2cset
/usr/bin/xargs
/usr/bin/hd
/bin/gzip
