# rom/usr/lib/opkg/info/dnidnsmasq.list
/etc/dnsmasq.conf
/sbin/mul_pppoe_dns
/etc/init.d/dnsmasq
/usr/sbin/dns-hijack
/usr/sbin/dnsmasq
