# rom/usr/lib/opkg/info/kmod-dni-mtdoops.list
/sbin/debug_wifi.sh
/etc/init.d/dni-debug.init
/sbin/wireless_log.sh
/sbin/console_log.sh
/sbin/mirror.sh
/sbin/capture_packet.sh
/sbin/arlo_log.sh
/sbin/wireless_client_statistic_log.sh
/lib/modules/3.14.77/mtdoops.ko
/sbin/thermal_log.sh
/sbin/debug_save_panic_log
/sbin/logread_log.sh
/sbin/basic_log.sh
/sbin/debug_here_log.sh
/sbin/enable_collect_console_log.sh
/sbin/save_collect_console_log.sh
/sbin/wlandebug.sh
/sbin/debug_circle.sh
/sbin/enable_ipv6_ping.sh
/sbin/debug_dutinfo.sh
/sbin/debug_telnetenable.sh
/sbin/debug_jumbo_frame.sh
