# rom/usr/lib/opkg/info/net-scan.list
/etc/init.d/satellite_ctype
/usr/sbin/net-scan
/usr/sbin/UA_Parser
/etc/device_type_table.json
/usr/sbin/dev-scan
/etc/init.d/netscan_if.init
/etc/init.d/netscan
/etc/oui36.csv
/etc/send_soap
/usr/sbin/devices_info
/etc/mam.csv
/etc/oui.csv
/etc/UA_String_Match_Table.txt
