# rom/usr/lib/opkg/info/openvpn_dni-openssl.list
/lib/upgrade/keep.d/openvpn_dni-openssl
/etc/openvpn/push_routing_rule
/etc/openvpn/vpn-firewall.sh
/etc/init.d/openvpn
/usr/sbin/openvpn
/etc/rc.d/S60openvpn
/etc/openvpn/download
/lib/upgrade/keep.d/openvpn
