# rom/usr/lib/opkg/info/qca-hostap_dni-10.4.list
/etc/init.d/qca-hostapd
/usr/sbin/hostapd
