# rom/usr/lib/opkg/info/qca-wsplcd_dni.list
/etc/hotplug.d/button/40-wsplcd
/usr/sbin/wsplcd
/etc/init.d/wsplcd
/etc/config/wsplcd
/sbin/uci2cfg.sh
/usr/sbin/ucpkgen
