# rom/usr/lib/opkg/info/qca-wpa-supplicant_dni-10.4.list
/etc/init.d/qca-wpa-supplicant
/usr/sbin/wpa_supplicant
