# rom/usr/lib/opkg/info/remGenie.list
/www/cgi-bin/genie.cgi
/usr/bin/genie_handler
/etc/init.d/remGenie
