# rom/usr/lib/opkg/info/udhcp-ext.list
/usr/share/udhcpc/default-arlo.script-ext
/usr/share/udhcpc/default.script.backhaul
/etc/init.d/setacl
/usr/share/udhcpc/default.script-ext
/sbin/udhcpc-ext
/sbin/udhcpd-ext
/usr/share/udhcpc/default.script.ap-ext
