# rom/usr/lib/opkg/info/wlan-common.list
/usr/sbin/post_wlan_updown.sh
/usr/sbin/wifison-sanity-check.sh
/etc/hotplug.d/wps/00-wps
/sbin/thermal_RA.sh
/sbin/update-wifi
/etc/hotplug.d/wlanlog/00-wlanlog
/sbin/factory-test.sh
/sbin/ap-led
/sbin/wlandebug
/lib/wifi/dniconfig.sh
/etc/hotplug.d/wps/01-wps-led
/etc/wifi_config_comp.awk
/etc/dni-wifi-config
/sbin/wpsled
/etc/init.d/wifison
/etc/init.d/wlan-common
/etc/ath/wifi.conf
/sbin/wifison.sh
/usr/sbin/check_psk_change_status
/usr/sbin/sync_button_led
/etc/search-wifi-interfaces.awk
/sbin/wlan_restart_queue.sh
/sbin/uci2dnicfg.sh
/sbin/wifimode-test.sh
/sbin/wlan
/sbin/wifison-monitor.sh
/usr/sbin/hw-revision-sanitycheck.sh
/sbin/wigig
/etc/ath/non_limit_country.conf
/etc/init.d/wifison-monitor
/usr/sbin/wireless_factory_config
