# usr/lib/opkg/info/auto-detection-tz.list
/usr/sbin/auto-detection-tz
/etc/init.d/get_xCloud_MAXmind.sh
