# usr/lib/opkg/info/base-files_dni.list
/lib/cfgmgr/cfgmgr.sh
/sbin/led_ring_action
/etc/init.d/done
/sbin/button_detect
/sbin/cfgapply
/etc/init.d/skb_recycler
/lib/cfgmgr/opmode.sh
/etc/hotplug2-init.rules
/sbin/check_status.sh
/etc/config/tt3
/sbin/wget2nand
/lib/upgrade/platform.sh
/sbin/cmd_sipalg
/etc/hotplug2-common.rules
/sbin/thermal.sh
/sbin/arlo_if_restart
/sbin/ping-netgear
/etc/hotplug.d/button/00-button
/etc/config/system
/sbin/vlan-dhcpd
/sbin/cfgmgr/cfg_email_ntp
/etc/init.d/powerctl
/etc/sysctl.conf
/sbin/printhosts
/etc/config/nowan
/sbin/setmac
/etc/init.d/detbutton
/sbin/bit_flip_check
/etc/init.d/ddns
/etc/init.d/opmode
/sbin/init-gpio
/etc/init.d/watchdog
/lib/preinit/80_mount_root
/etc/ntgrdata.cfg
/etc/init.d/defconfig
/lib/read_caldata_to_fs.sh
/sbin/reset_to_factory_settings.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/init.d/usb
/sbin/check_fwupgrade
/sbin/cfgmgr/cfg_apmode
/lib/cfgmgr/enet.sh
/sbin/cfgmgr/cfg_vlan_iptv
/sbin/cloud
/sbin/eth_handle
/lib/update_smp_affinity.sh
/lib/preinit/03_preinit_do_ipq806x.sh
/sbin/cmd_autoreset
/etc/ntgrdata-rbs20.cfg
/bin/login.sh
/etc/config/lanwan
/sbin/cfgmgr/cfg_opmode
/lib/update_system_params.sh
/etc/hotplug.d/usb/10-usb
/lib/ipq806x.sh
/sbin/cfgmgr/cfg_extmode
/etc/inittab
/sbin/drop_caches
/sbin/qos.sh
/sbin/wifi
/etc/init.d/enet
/sbin/kexec_load.sh
/sbin/run-ramfs
/sbin/cmd_ddns
/etc/hotplug.d/ieee1394/10-ieee1394
/etc/init.d/arlo
/etc/init.d/boot
/sbin/led_ring
/lib/usfd/ufsd.ko
/sbin/cfgmgr/cfg_brmode
/etc/uci-defaults/skb_recycler
/etc/uci-defaults/network
