# usr/lib/opkg/info/bluetopia.list
/usr/bin/libBTPSFILE.a
/usr/lib/libSS1BTVEN.a
/usr/bin/LinuxHCI
/usr/bin/LinuxGATM_CLT
/usr/bin/libBTPSVEND.a
/usr/bin/libSS1BTPM.a
/usr/bin/LinuxSPPLE
/usr/bin/LinuxSPPLE_PM
/usr/bin/LinuxL2CAP
/usr/bin/SS1BTPM
/usr/bin/libSS1BTGAT.a
/usr/bin/PS_KEY_CSR8811.txt
/usr/bin/libSS1BTGAT_LE.a
/usr/bin/libSS1BTPS.a
/usr/bin/ntgbpm
/usr/bin/LinuxSCO
/usr/bin/LinuxSDP
/etc/init.d/bluetopia
