# usr/lib/opkg/info/boxlogin.list
/usr/sbin/boxlogin
/etc/init.d/boxlogin
