# usr/lib/opkg/info/d2.list
/usr/include/d2lib/d2_pdu.h
/etc/init.d/d2_countryname.sh
/etc/d2/UpCfg.000.json
/etc/d2/XagentCfg.000.json
/usr/include/d2lib/sortedset.h
/usr/lib/libd2lib.so
/usr/include/d2lib/d2api.h
/etc/d2/RGenieCfg.000.json
/usr/bin/d2d
/usr/bin/d2
/etc/d2/ArloCfg.000.json
/etc/init.d/d2
