# usr/lib/opkg/info/dniminiupnpd.list
/usr/sbin/miniupnpd
/sbin/cmdupnp
/etc/init.d/upnp
