# usr/lib/opkg/info/kmod-qca-rfs_dni.list
/etc/modules.d/29-qca-rfs_dni
/lib/modules/3.14.77/qrfs.ko
/etc/init.d/qrfs
