# usr/lib/opkg/info/kmod-qca-hyfi-bridge_dni.list
/etc/modules.d/40-qca-hyfi-bridge_dni
/etc/init.d/hyfi-bridging
/lib/modules/3.14.77/hyfi-bridging.ko
