# usr/lib/opkg/info/kmod-atheros-gpio.list
/lib/modules/3.14.77/atherosgpio.ko
/etc/init.d/init-gpio
