# usr/lib/opkg/info/kmod-qca-wifi_dni-10.4-unified-profile.list
/lib/modules/3.14.77/umac.ko
/lib/modules/3.14.77/asf.ko
/lib/modules/3.14.77/ath_hal.ko
/usr/sbin/thermaltool
/usr/sbin/exttool
/usr/sbin/athadhoc
/etc/init.d/ssid_steering
/usr/sbin/apstats
/usr/sbin/spectraltool
/lib/modules/3.14.77/ath_rate_atheros.ko
/etc/config/ssid-steering
/lib/modules/3.14.77/smart_antenna.ko
/lib/modules/3.14.77/mem_manager.ko
/usr/sbin/assocdenialnotify
/usr/sbin/pktlogdump
/usr/sbin/80211stats
/usr/sbin/radartool
/lib/wifi/qcawifi_countrycode.txt
/lib/modules/3.14.77/qca_da.ko
/usr/sbin/wifitool
/usr/lib/libtlvparser.so
/usr/lib/libtlvencoder.so
/usr/sbin/athtestcmd
/usr/sbin/tx99tool
/lib/modules/3.14.77/ath_spectral.ko
/usr/sbin/athstats
/lib/wifi/33-qca-wifi_dni-10.4-unified-profile
/lib/modules/3.14.77/qdf.ko
/usr/lib/libtlvtemplate.so
/usr/lib/libtlvcmdrsp.so
/lib/modules/3.14.77/ath_dfs.ko
/lib/sbin/wifi_try
/lib/modules/3.14.77/wlan_ext.ko
/usr/lib/libtlvutil.so
/usr/sbin/wlanconfig
/lib/modules/3.14.77/ath_pktlog.ko
/usr/sbin/athdiag
/lib/modules/3.14.77/ath_dev.ko
/usr/sbin/athstatsclr
/usr/sbin/ssidsteering
/usr/sbin/wps_enhc
/usr/sbin/pktlogconf
/lib/modules/3.14.77/hst_tx99.ko
/lib/modules/3.14.77/qca_ol.ko
/usr/sbin/wifi_try
