# usr/lib/opkg/info/ntgdal.list
/usr/include/libdal.h
/usr/lib/libdal.so
/usr/bin/dal
/etc/init.d/ntgdal
/usr/bin/btdm
