# usr/lib/opkg/info/telnetenable.list
/lib/network/RtDebug.sh
/usr/sbin/telnetenable
/etc/init.d/telnet
