# usr/lib/opkg/info/tftpd-hpa.list
/etc/init.d/tftpd-hpa
/usr/sbin/tftpd-hpa
