# usr/lib/opkg/info/xagent.list
/etc/xagent/xagent/xagent.prod.env
/etc/xagent/remGenie/remGenie.qa.env
/etc/xagent/xagent/switch.sh
/etc/xagent/upagent/upagent.beta.env
/etc/xagent/upagent/upagent.qa2.env
/usr/bin/xagent
/etc/xagent/upagent/upagent.prod.env
/etc/xagent/remGenie/remGenie.qa2.env
/etc/xagent/remGenie/switch.sh
/etc/xagent/upagent/upagent.ftrial.env
/etc/xagent/upagent/upagent.staging.env
/etc/xagent/xcenvsw.sh
/etc/xagent/vzdaemon/vzdaemon.dev.env
/etc/xagent/upagent/upagent.qa.env
/etc/xagent/remGenie/remGenie.staging.env
/etc/xagent/xagent/xagent.dev.env
/etc/xagent/xagent/xagent.ftrial.env
/etc/xagent/vzdaemon/switch.sh
/etc/xagent/upagent/switch.sh
/etc/xagent/vzdaemon/vzdaemon.beta.env
/etc/xagent/vzdaemon/vzdaemon.qa2.env
/usr/share/ntgr/certs/ca-bundle-mega.crt
/etc/xagent/remGenie/remGenie.beta.env
/etc/xagent/xagent/xagent.staging.env
/usr/bin/evsinkd2
/etc/xagent/xagent/xagent.qa2.env
/etc/xagent/vzdaemon/vzdaemon.prod.env
/etc/xagent/upagent/upagent.dev.env
/etc/xagent/vzdaemon/vzdaemon.ftrial.env
/etc/xagent/remGenie/remGenie.prod.env
/usr/bin/xagent_control
/etc/xagent/vzdaemon/vzdaemon.staging.env
/etc/xagent/xagent/xagent.qa.env
/etc/xagent/remGenie/remGenie.ftrial.env
/etc/xagent/xagent/xagent.beta.env
/etc/xagent/remGenie/remGenie.dev.env
/etc/xagent/vzdaemon/vzdaemon.qa.env
/etc/init.d/XAgent
